# sshdfilter/Library/LaunchDaemons/net.jonbell.sshdfilter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.jonbell.sshdfilter</string>
	<key>ProgramArguments</key>
	<array>
		<string>/etc/sshdfilterLoad.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ServiceDescription</key>
	<string>SSHDFilter: load firewall rules</string>
</dict>
</plist>
